# browser_worker/deploy/browser-worker.service
[Unit]
Description=Searcher Browser Worker Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/searcher/browser_worker
EnvironmentFile=/opt/searcher/.env
ExecStart=/opt/searcher/browser_worker/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 8010 --root-path ${BROWSER_ROOT_PATH}
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
