# browser_worker/deploy/chromium-display.service
[Unit]
Description=Chromium Browser with Display and CDP
After=xvfb.service
Requires=xvfb.service

[Service]
Type=simple
User=root
Environment=DISPLAY=:99
ExecStart=/usr/bin/chromium \
    --no-sandbox \
    --disable-gpu \
    --remote-debugging-port=9222 \
    --remote-debugging-address=127.0.0.1 \
    --user-data-dir=/opt/searcher/browser_worker/chromium-profile \
    --disable-dev-shm-usage \
    --no-first-run \
    --no-default-browser-check \
    --window-size=1280,800
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
