# browser_worker/deploy/chromium-cdp.service
[Unit]
Description=Chromium CDP Remote Debugging Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/chromium \
    --headless=new \
    --no-sandbox \
    --disable-gpu \
    --remote-debugging-port=9222 \
    --remote-debugging-address=127.0.0.1 \
    --user-data-dir=/opt/repo/browser_worker/chromium-profile \
    --disable-dev-shm-usage \
    --no-first-run \
    --no-default-browser-check
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
