# browser_worker/deploy/x11vnc.service
[Unit]
Description=x11vnc VNC Server
After=xvfb.service
Requires=xvfb.service

[Service]
Type=simple
User=root
EnvironmentFile=/opt/searcher/.env
ExecStart=/usr/bin/x11vnc \
    -display :99 \
    -forever \
    -shared \
    -passwd ${VNC_PASSWORD} \
    -rfbport 5900 \
    -noxdamage
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
